# ledger.service
[Unit]
Description=Ledger Micronaut Native Service
# Wait for networking and DB to be ready before starting
After=network.target postgresql.service redis-server.service

[Service]
# Run as the restricted user
User=ledger
Group=ledger

# Directory where the app lives
WorkingDirectory=/opt/ledger

# 1. Load the Environment Variables automatically
EnvironmentFile=/opt/ledger/.env

# 2. Command to run the binary
ExecStart=/opt/ledger/ledger

# Restart automatically if it crashes
Restart=always
RestartSec=5

# Logging: Redirect Standard Output to System Logger
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ledger

[Install]
# Start this service when the system boots up
WantedBy=multi-user.target
